# deploy/scd-bot.service
[Unit]
Description=SCD WhatsApp Bot (FastAPI/uvicorn)
After=network.target

[Service]
Type=simple
User=scdbot
WorkingDirectory=/opt/scd-bot
EnvironmentFile=/opt/scd-bot/.env
# Bind to localhost only — Caddy terminates TLS and reverse-proxies to this.
ExecStart=/opt/scd-bot/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8000
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
